# solution-tp-final/backup_gitea.service
[Unit]
Description=Sauvegarde des données de Gitea

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup_gitea.sh