# solution-tp-final/cleanup_backups.service
[Unit]
Description=Nettoyage des anciennes sauvegardes de Gitea

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cleanup_backups.sh